# kbd-backlight-idle.service
[Unit]
Description=Keyboard Backlight Idle Timeout (Python/evdev)
# Ensure input-remapper is ready before we try to grab the device
After=input-remapper.service
# Stop the service if the target is no longer active
StopWhenUnneeded=yes

[Service]
Type=simple
# Standard practice to wait for hardware initialization
ExecStartPre=/usr/bin/sleep 5
# Call the python interpreter directly
ExecStart=/usr/bin/python3 /usr/local/bin/kbd-backlight-idle.py
Restart=always
RestartSec=5
# Ensure it runs as root to access /dev/input and /sys/class/leds
User=root

[Install]
WantedBy=multi-user.target suspend.target hibernate.target
